# .config/autostart/xfce4-clipman-plugin-autostart.desktop
[Desktop Entry]
Hidden=false
TryExec=sh -c "sleep 30; xfce4-clipman"
Exec=sh -c "sleep 30; xfce4-clipman"
